# roles/vdirsyncer/files/config/systemd/user/vdirsyncer.service
[Unit]
Description=Calendars and contacts synchronization service

[Service]
Type=oneshot
ExecStart=/bin/vdirsyncer sync

[Install]
WantedBy=default.target
